# myip.crontab
37 13 * * * cd $HOME/work/myip; ./myip.sh &> /tmp/myip.out
